# scripts/supervisor-zerotier.conf
[supervisord]
[program:zerotier]
command=zerotier-one
redirect_stdout=true
redirect_stderr=true